# ansible/roles/common/templates/gopro.service
[Unit]
Description=gopro uWSGI app
After=syslog.target

[Service]
ExecStart={{ install_dir }}/pyenv/bin/uwsgi \
        --socket :8000 --manage-script-name --mount /gopro=app:app \
	--virtualenv {{ install_dir }}/pyenv/ \
	--chdir {{ install_dir }}
User={{ user }}
Group=www-data
Restart=on-failure
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
